# deploy/starbot-webgui.service
# Systemd service file for Starbot WebGUI
# Place at: /etc/systemd/system/starbot-webgui.service
#
# Enable: sudo systemctl enable starbot-webgui
# Start:  sudo systemctl start starbot-webgui
# Status: sudo systemctl status starbot-webgui
# Logs:   sudo journalctl -u starbot-webgui -f

[Unit]
Description=Starbot WebGUI Service (Next.js)
After=network.target starbot-api.service

[Service]
Type=simple
User=starbot
Group=starbot
WorkingDirectory=/var/www/sites/stella/starbot.cloud
Environment="NODE_ENV=production"
Environment="PORT=3001"
Environment="NEXT_PUBLIC_API_URL=http://127.0.0.1:3737"
Environment="NEXT_PUBLIC_ADMIN_CONSOLE_HOSTS=console.sgauth0.com,www.console.sgauth0.com,console.starbot.cloud,www.console.starbot.cloud"
Environment="NEXT_PUBLIC_ADMIN_CONSOLE_URL=https://console.starbot.cloud"
Environment="NEXT_PUBLIC_CHAT_APP_URL=https://starbot.cloud"
Environment="ADMIN_CONSOLE_HOSTS=console.sgauth0.com,www.console.sgauth0.com,console.starbot.cloud,www.console.starbot.cloud"
Environment="ADMIN_EMAILS=stella@starbot.cloud"
Environment="ADMIN_CONSOLE_CODE=change-this-admin-code"
Environment="ADMIN_COOKIE_SECRET=change-this-cookie-secret"

# Start command (Next.js standalone server)
ExecStart=/usr/bin/node /var/www/sites/stella/starbot.cloud/.next/standalone/server.js

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=starbot-webgui

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
